# controller/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:ipcontroller]
command=/usr/local/bin/start_controller.sh
